# dark/luci-theme-noodlesDark-master/root/etc/uci-defaults/30_luci-theme-noodles
#!/bin/sh
uci batch <<-EOF
	set luci.themes.noodlesDark=/luci-static/noodles
	commit luci
EOF
exit 0
